# deployment/systemd/nexus-backend.service
[Unit]
Description=Nexus COS Node.js Backend
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/runner/work/nexus-cos/nexus-cos/backend
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npx ts-node src/server.ts
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexus-backend

[Install]
WantedBy=multi-user.target